# config/supervisord.conf
[supervisord]
nodaemon=true

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0777

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

[program:app]
command=/usr/local/bin/uwsgi --plugin python /etc/uwsgi.ini
stdout_logfile=/var/log/app.log
redirect_stderr=true
